# interactive/secretary-webhook.service
[Unit]
Description=LINE Interactive Secretary Webhook
After=network-online.target

[Service]
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/yuta/line/line-news-bot/interactive/run_server.sh
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
